# roles/python-services/templates/upstart.conf.j2
description "{{ service.name }} Service"

start on runlevel [2345]
stop on runlevel [!2345]

exec {{ service.install_path }}/app.py >> {{ service.install_path/log_file }} 2>&1